# configs/supervisord.conf
[program:pip-proxy]
command=/home/pi/pyenv/bin/python /srv/pip-proxy/pip-proxy.py
directory=/srv/pip-proxy/
user=pi
group=pi
environment=PATH="/home/pi/pyenv/bin"

